# exit-gate/python-app/exit-gate-gui.service
[Unit]
Description=Exit Gate GUI System
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=pi
Group=pi
Environment=DISPLAY=:0.0
Environment=HOME=/home/pi
WorkingDirectory=/home/pi/exit-gate/python-app
ExecStart=/bin/bash /home/pi/exit-gate/python-app/start_gui.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
